# com.obsidian-notify.watch-folders.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.obsidian-notify.watch-folders</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>YOUR_REPO_PATH/watch_folders.sh</string>
    </array>

    <!-- Runs at every 5-minute mark (StartCalendarInterval is more reliable than StartInterval for short-running scripts). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Minute</key><integer>0</integer></dict>
        <dict><key>Minute</key><integer>5</integer></dict>
        <dict><key>Minute</key><integer>10</integer></dict>
        <dict><key>Minute</key><integer>15</integer></dict>
        <dict><key>Minute</key><integer>20</integer></dict>
        <dict><key>Minute</key><integer>25</integer></dict>
        <dict><key>Minute</key><integer>30</integer></dict>
        <dict><key>Minute</key><integer>35</integer></dict>
        <dict><key>Minute</key><integer>40</integer></dict>
        <dict><key>Minute</key><integer>45</integer></dict>
        <dict><key>Minute</key><integer>50</integer></dict>
        <dict><key>Minute</key><integer>55</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>YOUR_REPO_PATH/notify.log</string>

    <key>StandardErrorPath</key>
    <string>YOUR_REPO_PATH/notify.log</string>
</dict>
</plist>
